# systemd/friend-report.service
[Unit]
Description=Daily friend (collab) activity report to ntfy

[Service]
Type=oneshot
ExecStart=/usr/local/bin/friend-report.sh
