# systemd/health-watchdog.service
[Unit]
Description=leaddneung health watchdog (ntfy alerts)
[Service]
Type=oneshot
ExecStart=%h/transcribe-opt/health-watchdog.sh
