# systemd/heartbeat.service
[Unit]
Description=leaddneung proof-of-life heartbeat to ntfy

[Service]
Type=oneshot
ExecStart=/home/neung/transcribe-opt/heartbeat.sh
