# systemd/heartbeat.timer
[Unit]
Description=Twice-daily leaddneung heartbeat (08:00 and 20:00)

[Timer]
OnCalendar=*-*-* 08:00:00
OnCalendar=*-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target
